# map3-edge/docker/geostat/rootfs/etc/services.d/geostat/run
#!/usr/bin/with-contenv sh

cd /opt/geostat && exec python /opt/geostat/geoparser.py

