# map3-edge/docker/geostat/rootfs/etc/services.d/geostat/run
#!/usr/bin/with-contenv sh

/optimize-influxdb.sh

exec /usr/bin/geostat

