# map3-edge-in-docker/rootfs/etc/services.d/docker-compose/run
#!/usr/bin/with-contenv sh

# pull images and run docker-compose up --no-build
exec /bash/docker-compose-up.sh

