# map3-edge-in-docker/rootfs/etc/services.d/docker-compose/run
#!/usr/bin/with-contenv sh

cd /src && exec /usr/bin/docker-compose up 

